# lab1/rootfs/var/spool/cron/crontabs/root
* * * * * /usr/local/ec535/print_time.sh &

